# scripts/com.opalrm.confluence-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.opalrm.confluence-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/douglasmaitland/Coding/rm-prototype/scripts/sync-from-confluence.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/douglasmaitland/Coding/rm-prototype</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/douglasmaitland/Coding/rm-prototype/scripts/confluence-sync.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/douglasmaitland/Coding/rm-prototype/scripts/confluence-sync.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
